# .config/systemd/user/led-sync.service
[Unit]
Description=Sync audio mute LEDs

[Service]
ExecStart=%h/.local/bin/led-sync.sh
Restart=always

[Install]
WantedBy=default.target