# .config/systemd/user/waybar-todo.service
[Unit]
Description=Refresh Waybar todo module

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill -RTMIN+12 waybar