# .config/systemd/user/waybar-todo.path
[Unit]
Description=Watch todo.txt for changes

[Path]
PathModified=%h/.todo/todo.txt
Unit=waybar-todo.service

[Install]
WantedBy=default.target